# init/init-redis.conf
description "start and stop redis-server"
version "1.0"

start on runlevel [2345]

script
exec /usr/local/bin/redis-server /etc/redis/redis.conf
end script
